# scheduler/systemd/hawkscapitol-secrets.timer
[Unit]
Description=Refresh HawksCapitol tmpfs secrets after boot

[Timer]
Unit=hawkscapitol-secrets.service
OnBootSec=30s
Persistent=false

[Install]
WantedBy=timers.target
