# scheduler/systemd/hawkscapitol-health-check.service
[Unit]
Description=HawksCapitol source and broker health check
Wants=network-online.target
Requires=hawkscapitol-secrets.service
After=network-online.target hawkscapitol-secrets.service

[Service]
Type=oneshot
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/HawksCapitol
Environment=HAWKSCAPITOL_ENV_FILE=/dev/shm/.hawkscapitol.env
Environment=HAWKSCAPITOL_REQUIRE_SHM=1
EnvironmentFile=/dev/shm/.hawkscapitol.env
ExecStart=/home/ec2-user/HawksCapitol/.venv/bin/python scheduler/run_health_check.py
