# launchd/com.hawkscapitol.daily-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hawkscapitol.daily-report</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HAWKSCAPITOL_HOME</key>
    <string>/path/to/HawksCapitol</string>
    <key>HAWKSCAPITOL_ENV_FILE</key>
    <string>/tmp/.hawkscapitol.env</string>
    <key>HAWKSCAPITOL_REQUIRE_SHM</key>
    <string>0</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "$HAWKSCAPITOL_HOME" &amp;&amp; set -a &amp;&amp; . "$HAWKSCAPITOL_ENV_FILE" &amp;&amp; set +a &amp;&amp; python3 scheduler/run_report.py</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/hawkscapitol-daily-report.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/hawkscapitol-daily-report.err</string>
</dict>
</plist>
